# login-monitor-1.1kord/debian/postinst
#!/bin/bash
systemctl enable login-monitor.service
