# login-monitor-1.1kord/debian/prerm
systemctl disable login-monitor.service
